# dpkg-skel/ccsyslogingester_VERSION-0ubuntu_amd64/lib/systemd/system/ccsyslogingester.service
[Unit]
Description=command center syslog ingester service
After=network.target

[Service]
Restart=always
RestartSec=1
User=ccsyslog
ExecStart=/usr/bin/ccsyslogingester --config=/etc/ccsyslog/config.json

[Install]
WantedBy=multi-user.target